# 008/vm/guest/etc/systemd/vagrant/log-monitor.timer
[Unit]
Description=Grep QUERY from FILE scheduled
Requires=log-monitor.service

[Timer]
Unit=log-monitor.service
OnUnitActiveSec=30

[Install]
WantedBy=timers.target