# 008/vm/guest/etc/systemd/vagrant/spawn-fcgi.service
[Unit]
Description=Spawn-fCgi Unit
After=network.target

[Service]
Type=simple
PIDFile=/var/run/spawn-fcgi.pid
EnvironmentFile=/etc/sysconfig/spawn-fcgi
ExecStart=/usr/bin/spawn-fcgi -n $OPTIONS
KillMode=process

[Install]
WantedBy=multi-user.target
