# 040/ansible/roles/barman_on_backup/files/services/barman-cron.service
[Unit]
Description=Barman cron
Wants=barman-cron.timer

[Service]
User=barman
#User=root
Type=oneshot
#Type=simple
ExecStart=/usr/bin/barman cron

[Install]
WantedBy=multi-user.target