# 040/ansible/roles/barman_on_backup/files/services/barman-cron.timer
[Unit]
Description=Barman cron scheduled
Requires=barman-cron.service

[Timer]
Unit=barman-cron.service
OnUnitActiveSec=10

[Install]
WantedBy=timers.target